# project/customfeed/MiWire/src/miwire.___init___
#!/bin/sh /etc/rc.common

START=99

USE_PROCD=1
PROG=/bin/miwire
CONFFILE=/etc/miwire/miwire.conf

start_service(){
    procd_open_instance
    procd_set_param command $PROG -C $CONFFILE
    procd_set_param file $CONFFILE
    procd_set_param respawn
    procd_append_param respawn 3600
	procd_append_param respawn 5
	procd_append_param respawn 10
    procd_set_param stdout 1 # forward stdout of the command to logd
    procd_set_param stderr 1 # same for stderr
    procd_close_instance
}
